# show_ip.service
[Unit]
Description=Show IP 
After=network.target

[Service]
Type=forking
ExecStart=/usr/bin/show_ip.sh

[Install]
WantedBy=multi-user.target
